# monorepo/cdk/shasta-cdk/remote.dev.loadtest/redpanda-bidirectional-s3-sync.service
[Unit]
Description=RedPanda Load Test Bidirectional S3 Sync with Hash Verification
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=exec
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/remote.dev.loadtest
ExecStart=/home/ec2-user/remote.dev.loadtest/bidirectional-s3-sync -upload-only -cleanup -hash sha256
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Environment variables for bidirectional sync
Environment=S3_BUCKET=redpanda-load-test-358474168551-us-east-1
Environment=LOGS_DIR=./logs
Environment=AWS_REGION=us-east-1
Environment=HASH_TYPE=sha256
Environment=CLEANUP_LOCAL=true
Environment=UPLOAD_ONLY=true

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/ec2-user/remote.dev.loadtest

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redpanda-bidirectional-s3-sync

# Resource limits
MemoryMax=512M
TasksMax=20

[Install]
WantedBy=multi-user.target 